# mirrulations-worker.service
[Unit]
Description=Mirrulations Download Worker
After=network.target redis6.service postgresql.service
Wants=redis6.service

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/mirrulations-search

# Environment — override these in /etc/systemd/system/mirrulations-worker.service.d/env.conf
# or set them directly here for the EC2 instance
EnvironmentFile=-/home/ec2-user/mirrulations-search/.env

Environment=REDIS_HOST=localhost
Environment=REDIS_PORT=6379
Environment=REDIS_DB=0
Environment=FETCH_REPO_DIR=/home/ec2-user/mirrulations-fetch
Environment=CSV_REPO_DIR=/home/ec2-user/mirrulations-csv

ExecStart=/home/ec2-user/mirrulations-search/.venv/bin/python /home/ec2-user/mirrulations-search/worker.py

# Restart policy: always restart on exit (crash or clean), with a short delay
Restart=always
RestartSec=5

# Kill cleanly on stop
KillSignal=SIGINT
TimeoutStopSec=30

StandardOutput=journal
StandardError=journal
SyslogIdentifier=mirrulations-worker

[Install]
WantedBy=multi-user.target
